# aws/system-update.cron
#!/bin/bash
apt-get update
apt-get dist-upgrade -y --no-install-recommends
shutdown -r now
